# src/server/main/services/Alarm.service
[Translation]
File=QtopiaServices
Context=Alarm
[Service]
Actions = "editAlarm()"
Icon = alarmbell
Name[]=Edit Alarm
