# src/server/main/services/Alert.service
[Translation]
File=QtopiaServices
Context=Alert
[Service]
Actions = "soundAlert()"
Icon = alarmbell
Name[]=Sound Alarm
