# src/server/main/services/calibrate.service
[Translation]
File=QtopiaServices
Context=Service:Name
[Service]
Icon=calibrate/Calibrate
Name[]=Calibrate Touchscreen
