# src/server/main/services/CallNetworks.service
[Translation]
File=QtopiaServices
Context=CallNetworks
[Service]
Icon=phonenetworks/PhoneNetworks
Name[]=Call Networks
