# src/server/main/services/CallHistory.service
[Translation]
File=QtopiaServices
Context=CallHistory
[Service]
Actions="showCallHistory(QCallList::ListType,QString);viewDetails(QCallListItem,QContact,int)"
Icon=callhistory/CallHistory
Name[]="Call History Details"
