# src/server/main/services/CallForwarding.service
[Translation]
File=QtopiaServices
Context=CallForarding
[Service]
Icon = callforwarding/CallForwarding
Name[]=Call Forwarding
