# src/server/main/services/CleanupWizard.service
[Translation]
File=QtopiaServices
Context=SysInfo
[Service]
Actions = "showCleanupWizard()"
Icon = sysinfo/SystemInfo
Name[]=Cleanup System
