# src/server/main/services/Clock.service
[Translation]
File=QtopiaServices
Context=Clock
[Service]
Actions = "showClock()"
Icon=clock/Clock
Name[]=Show Clock
