# src/server/main/services/Date.service
[Translation]
File=QtopiaServices
Context=Date
[Service]
Actions = "editDate()"
Icon = systemtime/DateTime
Name[]=Edit Date
