# src/server/main/services/Help.service
[Translation]
File=QtopiaServices
Context=Service:Name
[Service]
Actions="setDocument(QString)"
Icon=helpbrowser/HelpBrowser
