# src/server/main/services/InfraredBeaming.service
[Translation]
File=QtopiaServices
Context=InfraredBeaming
[Service]
Actions = "beamPersonalBusinessCard();beamBusinessCard(QContact);beamBusinessCard(QDSActionRequest);beamFile(QString,QString,QString,bool);beamCalendar(QDSActionRequest);beamFile(QContentId);"
Name[]=Send via Infrared
Icon=beam
[beamPersonalBusinessCard()]
Icon=beam
Name[]="Beam business card"
