# src/server/main/services/InstantMessage.service
[Translation]
File=QtopiaServices
Context=InstantMessage
[Service]
Actions = "writeMessage(QString)"
Icon = phone/sms
Name[]="Instant Message"
