# src/server/main/services/Launcher.service
[Translation]
File=QtopiaServices
Context=Service:Name
[Service]
Actions="execute(QString);execute(QString,QString)"
Icon=qpe/AppsIcon
