# src/server/main/services/MediaSession.service
[Translation]
File=QtopiaServices
Context=Service:Name
[Service]
Icon=multimedia
Name[]=Media Session
