# src/server/main/services/ModemEmulator.service
[Translation]
File=QtopiaServices
Context=Service:Name
[Service]
Name[]=Modem Emulator
Icon=qpe/ModemEmulator
