# src/server/main/services/NewEmailArrival.service
[Translation]
File=QtopiaServices
Context=Email
[Service]
Actions = "arrived(QDSActionRequest)"
Icon = email
Name[]=New Email Arrival
