# src/server/main/services/NewMmsArrival.service
[Translation]
File=QtopiaServices
Context=Mms
[Service]
Actions = "arrived(QDSActionRequest)"
Icon = phone/sms
Name[]=New MMS Arrival
