# src/server/main/services/NewSmsArrival.service
[Translation]
File=QtopiaServices
Context=Sms
[Service]
Actions = "arrived(QDSActionRequest)"
Icon = phone/sms
Name[]=New SMS Arrival
