# src/server/main/services/NewSystemMessageArrival.service
[Translation]
File=QtopiaServices
Context=System
[Service]
Actions = "arrived(QDSActionRequest)"
Icon = phone/sms
Name[]=New System Message Arrival
