# src/server/main/services/PhotoEdit.service
[Translation]
File=QtopiaServices
Context=Photo Edit
[Service]
Actions = "editImage(QDSActionRequest)"
Icon=photoedit/ImageViewer
Name[]=Photo Edit
