# src/server/main/services/Presence.service
[Translation]
File=QtopiaServices
Context=Presence
[Service]
Actions = "dnd()"
Icon=presence-none
Name[]=Enable or disable Do not disturb status
