# src/server/main/services/QtopiaPowerManager.service
[Translation]
File=QtopiaServices
Context=Service:Name
[Service]
Icon=light-and-power/Light
Name[]=Qt Extended Power Manager
