# src/server/main/services/Receive.service
[Translation]
File=QtopiaServices
Context=Receive
[Service]
Actions = "receiveData(QString,QString)"
