# src/server/main/services/RotationManager.service
[Translation]
File=QtopiaServices
Context=Service:Name
[Service]
Icon=rotation/Rotation
Name[]=Rotation Manager
