# src/server/main/services/Settings.service
[Translation]
File=QtopiaServices
Context=CaptureSetting
[Service]
Multiple = 1
