# src/server/main/services/SettingsManager.service
[Translation]
File=QtopiaServices
Context=CaptureSetting
[Service]
Name[]=Settings Manager
Icon=settings
