# src/server/main/services/Suspend.service
[Translation]
File=QtopiaServices
Context=Service:Name
[Service]
Icon=light-and-power/Light
Name[]=Media Session
