# src/server/main/services/SystemMessages.service
[Translation]
File=QtopiaServices
Context=Service:Name
[Service]
Name[]=SystemMessages
Icon=phone/sms
