# src/server/main/services/Telephony.service
[Translation]
File=QtopiaServices
Context=Service:Name
[Service]
Multiple=1
