# src/server/main/services/TimeMonitor.service
[Translation]
File=QtopiaServices
Context=Service:Name
[Service]
Actions="timeChange(QString)"
Multiple=1
