# src/server/main/services/Time.service
[Translation]
File=QtopiaServices
Context=Time
[Service]
Actions = "editTime()"
Icon = systemtime/DateTime
Name[]=Edit Time
