# src/server/main/services/UsbGadget.service
[Translation]
File=QtopiaServices
Context=UsbGadget
[Service]
Name[]="USB Gadget"
Multiple=1
