# src/server/main/services/VirtualKeyboard.service
[Translation]
File=QtopiaServices
Context=Service:Name
[Service]
Icon=qpe/keyboard
Name[]="Virtual Keyboard"
