# src/server/main/services/VoiceMail.service
[Translation]
File=QtopiaServices
Context=VoiceMail
[Service]
Actions = "setVoiceMail()"
Icon=phone/voicemail
Name[]=Edit Voice Mail Number
