# src/server/main/services/UsbGadget/Ethernet.service
[Translation]
File=QtopiaServices
Context=UsbGadget/Ethernet
[Service]
Actions="startUsbService(QString)"
Icon=qpe/USBEthernet
Name[]="USB Ethernet Service"
Multiple=1
Gadget=QUsbEthernetGadget
