# src/server/main/services/UsbGadget/Serial.service
[Translation]
File=QtopiaServices
Context=UsbGadget/Serial
[Service]
Actions="startUsbService(QString)"
Icon=qpe/USBSerial
Name[]="USB Serial Service"
Gadget=QUsbSerialGadget
